# deb/etc/profile.d/99-xamarin-android.sh
#!/bin/bash
export PATH=$PATH:/opt/xammarin.android/bin/Release/bin
